# crontab
*/10 * * * * python /ks-notify.py >> /logfile